# screen.service
[Unit]
Description=screen

[Service]
Type=simple
ExecStart=/opt/screen.sh

[Install]
WantedBy=multi-user.target
